# devrel-integration/agentic-base-bot.service
[Unit]
Description=Agentic-Base Integration Bot
Documentation=https://github.com/your-org/agentic-base
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agentic-base
Group=agentic-base
WorkingDirectory=/opt/agentic-base/integration

# Environment file with secrets
EnvironmentFile=/opt/agentic-base/integration/secrets/.env.local

# Additional environment variables
Environment="NODE_ENV=production"
Environment="TZ=UTC"

# Main process
ExecStart=/usr/bin/node dist/bot.js

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Resource limits
# Limit memory usage to 512MB
MemoryMax=512M
# Limit CPU usage to 100%
CPUQuota=100%

# Security hardening
# Run as non-root user
NoNewPrivileges=true
# Restrict file system access
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agentic-base/integration/logs
ReadWritePaths=/opt/agentic-base/integration/data

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentic-base-bot

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
